# Slave/supervisord.conf
[supervisord]
nodaemon=true

[program:socket-master]
command=python socket_slave.py
autostart=true
autorestart=false